# ~/.local/share/autostart/agent-core-ng-final.desktop
[Desktop Entry]
Type=Application
Name=Agent Core NG Final Tray
Comment=System tray application for Agent Core NG
Exec=env DISPLAY=:0 XAUTHORITY=/home/art/.Xauthority /usr/bin/python3 /home/art/agent-RegArt-1/tray-app/final_working_tray.py
Icon=/home/art/agent-RegArt-1/tray-app/favicon.png
Terminal=false
X-GNOME-Autostart-enabled=true